# timers/waterholes-update.timer
[Unit]
Description=Update Waterholes for IHM every 10 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
